# defaults/scripts/adone
#!/sbin/runscript

extra_started_commands="reload"

ADONE=%ADONE_PATH%
USER=%USER%

export PATH=%NODE_PATH%:$PATH
export ADONE_HOME="%HOME_PATH%"

super() {
    su - $USER -c "PATH=$PATH; ADONE_HOME=$ADONE_HOME $*"
}

depend() {
	use net localmount
	after bootmisc
}

start() {
	ebegin "Starting ADONE"

	super $ADONE resurrect

	eend $?
}

stop() {
	ebegin "Stopping ADONE"

	#super $ADONE dump
	super $ADONE delete all
	super $ADONE kill

	eend
}

reload() {
	ebegin "Reloading ADONE"

	super $ADONE reload all

	eend $?
}

status() {
	ebegin "Status for ADONE"

	super $ADONE list

	eend $?
}
